# host/systemd/system/ollie_watchdog.timer
[Unit]
Description=Run Ollie Watchdog every minute

[Timer]
# 開機後延遲 1 分鐘啟動第一次檢查
OnBootSec=1min
# 之後每 1 分鐘執行一次
OnUnitActiveSec=1min
Unit=ollie_watchdog.service

[Install]
WantedBy=timers.target
